# sample/crontab
0 5 * * * /home/ubuntu/capstone/start.sh